# brotherscankeyd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting brotherscankeyd"
	/usr/bin/brotherscankeyd.py --syslog daemon
	eend $?
}

stop() {
	ebegin "Stopping brotherscankeyd"
	/usr/bin/brotherscankeyd.py stop
	eend $?
}
