# brotherscankeyd
#!/sbin/runscript

depend() {
	need net
}

start() {
	local verbose
	yesno $rc_verbose || verbose=--verbose
	ebegin "Starting brotherscankeyd"
	/usr/bin/brotherscankeyd.py ${verbose} --syslog daemon
	eend $?
}

stop() {
	ebegin "Stopping brotherscankeyd"
	/usr/bin/brotherscankeyd.py stop
	eend $?
}
